# config/launchd/ghost.mcp-proxy.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL_PREFIX__.mcp-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV__/bin/python3</string>
        <string>-m</string>
        <string>ghost.services.mcp_proxy</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__GHOST_HOME__/git/ghost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- WorkingDirectory + python -m adds the repo to sys.path automatically -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>MCP_PROXY_PORT</key>
        <string>__MCP_PROXY_PORT__</string>
        <key>MCP_BACKEND_PORT</key>
        <string>__MCP_BACKEND_PORT__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__GHOST_HOME__/ghost_run_dir/mcp-proxy.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__GHOST_HOME__/ghost_run_dir/mcp-proxy.stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
