# config/launchd/ghost.daemon.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL_PREFIX__.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__GHOST_HOME__/git/ghost/ghost/bin/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__GHOST_HOME__/git/ghost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- start.sh self-locates GHOST_HOME from its own path — no hardcoded paths here -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>MCP_PROXY_PORT</key>
        <string>__MCP_PROXY_PORT__</string>
        <key>MCP_BACKEND_PORT</key>
        <string>__MCP_BACKEND_PORT__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__GHOST_HOME__/ghost_run_dir/ghost.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__GHOST_HOME__/ghost_run_dir/ghost.stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
